# tomcat/files/tomcat.service.jinja
# This file was generated by SaltStack with tomcat-formula.
# Manual changes will be overwritten at the next run of SaltStack.
# Last run: {{ None|strftime("%c") }}
[Unit]
Description=Apache Tomcat
Documentation=http://tomcat.apache.org
Requires=network.target remote-fs.target
After=network.target remote-fs.target

[Service]
Type=forking

Environment='JAVA_HOME={{ tomcat_settings.env.java_home }}'
Environment='JRE_HOME={{ tomcat_settings.env.java_home }}/jre'
Environment='CATALINA_HOME=/opt/{{ tomcat_name }}'
Environment='CATALINA_BASE=/opt/{{ tomcat_name }}'
Environment='CATALINA_PID=/opt/{{ tomcat_name }}/temp/tomcat.pid'
Environment='CATALINA_OPTS=-Xms512M -Xmx1024M -server -XX:+UseParallelGC'

User={{ tomcat_settings.env.user }}
Group={{ tomcat_settings.env.group }}
UMask=0007
RestartSec=10
Restart=always

ExecStart=/opt/tomcat/bin/startup.sh
ExecStop=/opt/tomcat/bin/shutdown.sh

Restart=on-failure
SyslogIdentifier=tomcat
WorkingDirectory=/home/{{ tomcat_settings.env.user  }}

[Install]
WantedBy=multi-user.target
